# config/mora.socket
# -*- conf -*-

[Unit]
Description=MORA socket

[Socket]
ListenStream=/run/mora/socket

[Install]
WantedBy=sockets.target
